# arm_chromebook_debian/.config/autostart/pi-apps-updater.desktop
[Desktop Entry]
Name=Pi-Apps Updater
Exec=/home/linux/pi-apps/updater onboot
Icon=/home/linux/pi-apps/icons/logo.png
Terminal=false
StartupWMClass=Pi-Apps
Type=Application
X-GNOME-Autostart-enabled=true
Hidden=false
NoDisplay=false
